# tools/etc/bench_agent.ini
[program:isucon]
command=/Users/myfinder/Development/github/isucon/tools/etc/agent.sh
process_name=isucon bench agent
stdout_logfile=/tmp/isucon.agent.log
stderr_logfile=/tmp/isucon.agent.log
