# tools/etc/bench_master.ini
[program:isucon]
command=/Users/myfinder/Development/github/isucon/tools/etc/master.sh
process_name=isucon bench master
stdout_logfile=/tmp/isucon.master.log
stderr_logfile=/tmp/isucon.master.log
